# supervisord.conf
[supervisord]
logfile=/var/log/supervisord.log ; (main log file;default $CWD/supervisord.log)
nodaemon=true
childlogdir=/var/log/supervisord

[program:dynamodb]
command=java -Djava.library.path=/usr/lib/DynamoDBLocal_lib
  -jar /usr/lib/DynamoDBLocal.jar
  -port 8000
  -sharedDb
  %(ENV_DYNAMODB_LOCAL_ARGS)s

[program:dynamo-admin]
command=/usr/local/bin/dynamodb-admin -p 8080

[program:nginx]
command=/usr/sbin/nginx -c /etc/nginx-proxy.conf
